# Raspberry_Pi/etc/systemd/system/wifi-fallback.service
[Unit]
Description=WiFi Fallback Script
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wifi-fallback.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
